# tools/crontab
0 4 0 0 */2 root /ruby/urbanevents/tools/cronall.sh >> /var/log/cron.log 2>&1
