# tools/crontab
* * * * * root /ruby/urbanevents/tools/cronrunner.sh output.sh
0 10 0 0 */2 root /ruby/urbanevents/tools/cronrunner.sh backup_es.sh
